# deploy/easy-week-backend.service
[Unit]
Description=Easy Week backend (FastAPI + Uvicorn)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# ВАЖНО: поправь путь и пользователя под свой Raspberry Pi
User=pashtitto
WorkingDirectory=/home/pashtitto/easy-week/backend
EnvironmentFile=/home/pashtitto/easy-week/backend/.env
# Порт 8010 — свободный (perdulário занимает 3001/5432/80, «другой проект» — 3000)
ExecStart=/home/pashtitto/easy-week/backend/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8010
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
